# kafka/supervisor/kafka.conf
[program:kafka]
command=/usr/local/bin/start-kafka.sh
autostart=true
autorestart=true
