# kafka/supervisor/ocsp.conf
[program:ocsp]
command=/usr/bin/start-ocsp.sh
autostart=true
autorestart=true