# kafka/supervisor/zookeeper.conf
[program:zookeeper]
command=/usr/bin/start-zookeeper.sh
autostart=true
autorestart=true
