# kafkaproxy/kafkaproxy.conf
[program:kafka_proxy]
command=/start_proxy.sh
autostart=true
autorestart=true
