# kafkacluster/supervisor/kafka2.conf
[program:kafka2]
command=/usr/bin/start-kafka.sh 2 9093
autostart=true
autorestart=true
