# kafkacluster/supervisor/kafka1.conf
[program:kafka1]
command=/usr/bin/start-kafka.sh 1 9092
autostart=true
autorestart=true
